# server-setup/shiboroom-backend.service
[Unit]
Description=Shiboroom Backend API Server
After=network.target mysql.service meilisearch.service

[Service]
Type=simple
User=grik
Group=grik
WorkingDirectory=/var/www/shiboroom/backend
ExecStart=/var/www/shiboroom/backend/shiboroom-api

# 環境変数
Environment="CONFIG_PATH=/var/www/shiboroom/config/scraper_config.yaml"
Environment="PORT=8085"
Environment="ENV=production"

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shiboroom-backend

# 再起動設定
Restart=always
RestartSec=5

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
